# server-rust-luc/root/etc/init.d/server-rust-luc
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: GPL-2.0-only
# Copyright (C) 2024-2026 yorkane

USE_PROCD=1
START=95
STOP=15

NAME="server-rust-luc"
CONF="server-rust-luc"
CONFDIR="/var/etc/${NAME}"

generate_server_config() {
	local section="$1"
	local enabled server server_port password method timeout no_delay ipv6_first mode plugin plugin_opts acl

	config_get_bool enabled "$section" 'enabled' 0
	[ "$enabled" = "1" ] || return 0

	config_get server "$section" 'server' '0.0.0.0'
	config_get server_port "$section" 'server_port' '8388'
	config_get password "$section" 'password' ''
	config_get method "$section" 'method' 'aes-256-gcm'
	config_get timeout "$section" 'timeout' '300'
	config_get_bool no_delay "$section" 'no_delay' 0
	config_get_bool ipv6_first "$section" 'ipv6_first' 0
	config_get mode "$section" 'mode' 'tcp_and_udp'
	config_get plugin "$section" 'plugin' ''
	config_get plugin_opts "$section" 'plugin_opts' ''
	config_get acl "$section" 'acl' ''

	[ -z "$password" ] && {
		echo "server-rust-luc: skipping server $section - no password set" >&2
		return 0
	}

	local config_file="${CONFDIR}/${section}.json"

	json_init
	json_add_string "server" "$server"
	json_add_int "server_port" "$server_port"
	json_add_string "password" "$password"
	json_add_string "method" "$method"
	json_add_int "timeout" "$timeout"

	[ "$no_delay" = "1" ] && json_add_boolean "no_delay" 1
	[ "$ipv6_first" = "1" ] && json_add_boolean "ipv6_first" 1
	json_add_string "mode" "$mode"

	if [ -n "$plugin" ]; then
		json_add_string "plugin" "$plugin"
		[ -n "$plugin_opts" ] && json_add_string "plugin_opts" "$plugin_opts"
	fi

	json_dump > "$config_file"
	json_cleanup

	procd_open_instance "server-${section}"
	procd_set_param command /usr/bin/ssserver -c "$config_file"
	[ -n "$acl" ] && procd_append_param command --acl "$acl"
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param file "$config_file"
	procd_close_instance

	echo "server-rust-luc: started server instance $section on port $server_port" >&2
}

generate_manager_config() {
	local section="$1"
	local enabled manager_address manager_port method timeout mode

	config_get_bool enabled "$section" 'enabled' 0
	[ "$enabled" = "1" ] || return 0

	config_get manager_address "$section" 'manager_address' '127.0.0.1'
	config_get manager_port "$section" 'manager_port' '6100'
	config_get method "$section" 'method' 'aes-256-gcm'
	config_get timeout "$section" 'timeout' '300'
	config_get mode "$section" 'mode' 'tcp_and_udp'

	local config_file="${CONFDIR}/manager_${section}.json"
	local manager_addr="${manager_address}:${manager_port}"

	# Build JSON config for manager with user list
	local servers_json=""
	local first=1

	generate_manager_user() {
		local user_section="$1"
		local user_port user_password user_method

		config_get user_port "$user_section" 'server_port'
		config_get user_password "$user_section" 'password'
		config_get user_method "$user_section" 'method' "$method"

		[ -z "$user_port" ] && return 0
		[ -z "$user_password" ] && return 0

		if [ "$first" = "1" ]; then
			first=0
		else
			servers_json="${servers_json},"
		fi

		servers_json="${servers_json}{\"server_port\":${user_port},\"password\":\"${user_password}\",\"method\":\"${user_method}\"}"
	}

	config_foreach generate_manager_user manager_user

	cat > "$config_file" <<-EOF
	{
	    "manager_address": "${manager_addr}",
	    "servers": [${servers_json}],
	    "method": "${method}",
	    "timeout": ${timeout},
	    "mode": "${mode}"
	}
	EOF

	procd_open_instance "manager-${section}"
	procd_set_param command /usr/bin/ssmanager -c "$config_file"
	procd_set_param respawn
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_set_param file "$config_file"
	procd_close_instance

	echo "server-rust-luc: started manager instance $section" >&2
}

start_service() {
	local main_enabled
	config_load "${CONF}"
	config_get_bool main_enabled main enabled 0

	[ "$main_enabled" = "1" ] || {
		echo "server-rust-luc: service disabled (set enabled=1 in general section)" >&2
		return 0
	}

	mkdir -p "${CONFDIR}"

	. /usr/share/libubox/jshn.sh

	config_foreach generate_server_config server
	config_foreach generate_manager_config manager
}

stop_service() {
	rm -rf "${CONFDIR}"
}

reload_service() {
	stop
	start
}

service_triggers() {
	procd_add_reload_trigger "${CONF}"
}
